# profile.d/ixon.sh
#
# Recovering from Ctrl+S in terminal
#

# Enable XON/XOFF flow control
stty -ixon
